# templates/default/kibana4.sh.erb
#!/bin/bash
#
# chkconfig: 35 90 12
# description: kibana4 server

#configure this with wherever you unpacked kibana:
KIBANA_BIN=<%= node['kibana']['installdir'] %>/bin

PID_FILE=/var/run/$NAME.pid
PATH=/bin:/usr/bin:/sbin:/usr/sbin:$KIBANA_BIN
DAEMON=$KIBANA_BIN/kibana
NAME=kibana4
DESC="Kibana4"


#
# Get function from functions library
. /etc/init.d/functions
# Start the service kibana4
start() {
        initlog -c "echo -n Starting kibana4 server: "
        $DAEMON &
        ### Create the lock file ###
        touch /var/lock/subsys/kibana4
        success $"kibana4 server startup"
        echo
}
# Restart the service FOO
stop() {
        initlog -c "echo -n Stopping kibana4 server: "
        killproc kibana4
        ### Now, delete the lock file ###
        rm -f /var/lock/subsys/kibana4
        echo
}
### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status kibana4
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit 0



